# USBsentry/usb-scan@.service
# =============================================================================
# usb-scan@.service
# systemd service triggered by udev when a USB device is inserted
# =============================================================================
#
# INSTALLATION:
#   sudo cp usb-scan@.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#
# The @ symbol makes this a template service. The %i specifier is replaced
# with the device name passed from the udev rule (e.g., sda1).
#
# =============================================================================

[Unit]
Description=ClamAV USB Security Scan for /dev/%i
After=clamav-daemon.service
Requires=clamav-daemon.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/usb_scan.sh /dev/%i
TimeoutStartSec=360
StandardOutput=journal
StandardError=journal

# Run as root — required for raw device access and mount operations
User=root
Group=root
